# cron/sysstat.cron.hourly.in
#!/bin/sh
# Run system activity accounting tool every @COLLECT_INTERVAL@ minute(s)
@SA_LIB_DIR@/sa1 @COLLECT_INTERVAL_SEC@ @COLLECT_COUNT@
