# cron/sysstat-collect.timer.in
# @SYSTEMD_UNIT_DIR@/sysstat-collect.timer
# (C) 2014 Tomasz Torcz <tomek@pipebreaker.pl>
#
# @PACKAGE_NAME@-@PACKAGE_VERSION@ systemd unit file:
#        Activates activity collector every @COLLECT_INTERVAL@ minute(s)

[Unit]
Description=Run system activity accounting tool every @COLLECT_INTERVAL@ minute(s)

[Timer]
OnCalendar=*:00/@COLLECT_INTERVAL@

[Install]
WantedBy=sysstat.service
