# cron/sysstat-rotate.timer.in
# @SYSTEMD_UNIT_DIR@/sysstat-rotate.timer
# (C) 2023 Sebastien Godard <sysstat_at_orange_dot_fr>
#
# @PACKAGE_NAME@-@PACKAGE_VERSION@ systemd unit file:
#        Rotate daily system activity data file at midnight

[Unit]
Description=Rotate daily system activity data file at midnight

[Timer]
OnCalendar=00:00

[Install]
WantedBy=sysstat.service
